# ota-server/ota-server.service
[Unit]
Description=OTA Firmware Management Server
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/ota-server
ExecStart=/opt/ota-server/venv/bin/gunicorn -w 2 -b 127.0.0.1:5000 app:app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
